# converted_sigma/linux/file_event/file_event_lnx_persistence_cron_files_esql.yml
rule_id: 943dc353-28e7-40d6-af5c-dced52a5c3a2
name: "Persistence via Cron Files"
description: |
  Detects creation of cron files in cron directories.
  Adversaries may create cron jobs to establish persistence.
  Monitors /etc/cron.d, /etc/cron.daily, etc.
type: esql
params: '{}'
schedule_interval: 10m
index:
  - auditbeat-*
query_language: esql
query: |
  FROM auditbeat-*
  /* file creation events */
  | WHERE event.category == "file" AND event.type == "creation"
  /* selection: cron directories */
  | WHERE (
      file.path RLIKE "(?i)^/etc/cron\\.d/"
      OR file.path RLIKE "(?i)^/etc/cron\\.daily/"
      OR file.path RLIKE "(?i)^/etc/cron\\.hourly/"
      OR file.path RLIKE "(?i)^/etc/cron\\.monthly/"
      OR file.path RLIKE "(?i)^/etc/cron\\.weekly/"
      OR file.path RLIKE "(?i)^/var/spool/cron/crontabs/"
      OR file.path IN ("/etc/cron.allow", "/etc/cron.deny", "/etc/crontab")
    )
  /* Filter: exclude OS package managers and updates */
  | WHERE NOT (process.name IN ("dpkg", "yum", "dnf", "apt", "dockerd", "ansible-playbook"))
  | STATS
      event_count = COUNT(*),
      users = VALUES(user.name),
      processes = VALUES(process.name),
      file_paths = VALUES(file.path),
      hosts = VALUES(host.name)
    BY host.name, user.name, bucket = DATE_TRUNC(20 minutes, @timestamp)
  | WHERE event_count >= 1
  | KEEP agent_id, bucket, host.name, user.name, event_count, users, processes, file_paths, hosts
  | SORT event_count DESC
enabled: false
tags:
  - attack.persistence
  - attack.privilege-escalation
  - attack.t1053.003
severity: medium
risk_score: 60
references:
  - "https://attack.mitre.org/techniques/T1053/003/"
nist: ["SI-4", "AU-12"]
pci_dss: ["10.2"]
gdpr: ["32.1.b"]
mitre_attack:
  TA0003:
    T1053: ["T1053.003"]
