# pci-dss-all/PCI DSS 11.5 - Critical File Modification Detected (Linux).yml
rule_id: f4da2572-0860-4357-88ad-d88ed2e83274
name: "PCI DSS 11.5 - Critical File Modification Detected (Linux)"
description: |
  Detects unauthorized modification of critical system files on Linux via auditd FIM rules.
  PCI DSS Requirement 11.5 requires a change-detection mechanism for critical system files,
  configuration files, and content files with alerts for unauthorized modifications.
  Monitors writes/attribute changes to critical paths: /etc/passwd, /etc/shadow,
  /etc/sudoers, /etc/ssh/*, /etc/cron*/*, /etc/hosts, /etc/pam.d/*, 
  /usr/bin/, /usr/sbin/, /bin/, /sbin/ etc.
  Requires auditd rules configured for -w (watch) on critical file paths.
  Trigger: Any write, attribute change on monitored critical file paths.
  False positives: Legitimate package updates (yum/apt), authorized admin changes.
  Triage: Identify the process and user that made the change, check against
  change management system, verify package manager was the actor vs manual edit.
type: esql
params: '{}'
schedule_interval: 5m
index:
  - logs-audit*
query_language: esql
query: |
  FROM logs-audit*
  | WHERE event.dataset == "auditd.log"
    AND event.action IN ("opened-file", "changed-file-attributes-of",
                          "write", "truncated-file", "renamed",
                          "changed-file-permissions-of")
    AND (
      file.path LIKE "/etc/passwd*"
      OR file.path LIKE "/etc/shadow*"
      OR file.path LIKE "/etc/sudoers*"
      OR file.path LIKE "/etc/ssh/*"
      OR file.path LIKE "/etc/cron*"
      OR file.path LIKE "/etc/pam.d/*"
      OR file.path LIKE "/etc/hosts*"
      OR file.path LIKE "/etc/ld.so*"
      OR file.path LIKE "/bin/*"
      OR file.path LIKE "/sbin/*"
      OR file.path LIKE "/usr/bin/*"
      OR file.path LIKE "/usr/sbin/*"
      OR file.path LIKE "/usr/lib/*"
      OR file.path LIKE "/boot/*"
    )
    AND NOT process.name IN ("yum", "apt", "apt-get", "dpkg", "rpm", "dnf", "zypper")
  | STATS
      file_changes       = COUNT(*),
      unique_files       = COUNT_DISTINCT(file.path),
      modified_files     = VALUES(file.path),
      actors             = VALUES(user.name),
      processes_used     = VALUES(process.name),
      change_types       = VALUES(event.action)
    BY host.name,
       user.name,
       bucket = DATE_TRUNC(5 minutes, @timestamp)
  | EVAL severity_level = CASE(
      unique_files >= 10,                    "critical",
      unique_files >= 5,                     "high",
      unique_files >= 2,                     "medium",
      "low"
    )
  | EVAL risk_score = CASE(
      unique_files >= 10, 95,
      unique_files >= 5,  75,
      unique_files >= 2,  55,
      40
    )
  | KEEP bucket, host.name, user.name, file_changes,
         unique_files, modified_files, actors, processes_used,
         change_types, severity_level, risk_score
  | SORT risk_score DESC
enabled: false
tags:
  - pci_dss
  - file_integrity
  - linux
  - fim
severity: medium
risk_score: 55
references:
  - "PCI DSS v3.2.1 Requirement 11.5: Deploy change-detection mechanism for critical files"
  - "PCI DSS v3.2.1 Requirement 11.5.1: Respond to alerts from change-detection solution"
pci_dss: ["11.5", "11.5.1"]
nist: ["CM-3","CM-6","SI-7","AU-12"]
gdpr: ["5.1.f", "25", "32"]
mitre_attack:
  TA0003:
    T1543: []
  TA0005:
    T1222: []
  TA0040:
    T1565: []