# categoried/Endpoint-OS/Scheduled Task or Cron Job Suspicious Modification.yml
rule_id: e6fb7e81-c30b-4ac6-b892-5c63fa329772
name: "Scheduled Task or Cron Job Suspicious Modification"
description: |
  This rule detects creation or modification of scheduled tasks (Windows Task Scheduler) or cron jobs that could establish persistence. For Windows, it monitors schtasks.exe commands creating tasks with suspicious patterns. For Linux, it detects crontab modifications or direct /etc/cron.* file edits. The rule triggers when 2+ scheduled task operations occur on a host within 10 minutes.

#  **Why it matters:**
#  Scheduled tasks are a primary persistence mechanism across all operating systems:
#  - MITRE ATT&CK T1053 documents scheduled task/job abuse for persistence and privilege escalation
#  - NIST SP 800-53 CM-3 requires change control for system configurations including schedulers
#  - PCI DSS 10.2.5 mandates logging of changes to system components
#
#  **False Positives:**
#  - Legitimate automation scripts creating scheduled tasks for business processes
#  - System administrators configuring backup or maintenance jobs
#  - Software installations registering update checks or telemetry tasks
#  - Configuration management systems (SCCM, GPO) deploying scheduled tasks
#
#  **Triage Steps:**
#  1. Review the scheduled task command and arguments for suspicious behavior
#  2. Identify the user who created the task and verify authorization
#  3. Check task scheduling parameters - does it run at suspicious times (e.g., startup, hourly)?
#  4. Examine the executable or script being scheduled for malicious indicators
#  5. Verify if the task was created through a change management process
#  6. Disable suspicious tasks and investigate the creating user's other activity
type: esql
params: '{}'
schedule_interval: 10m
index:
  - winlogbeat-*
  - auditbeat-*
query_language: esql
query: |
  FROM winlogbeat-*, auditbeat-*
  // Time window: last 10 minutes
  | WHERE @timestamp >= NOW() - 10 minutes
  // Process creation or file modification events
  | WHERE event.code == "4688" OR event.action RLIKE ".*(exec|created|modified).*"
  // Extract command line and file path
  | EVAL cmd = COALESCE(process.command_line, winlog.event_data.CommandLine, event.original)
  | EVAL filepath = COALESCE(file.path, "")
  | EVAL hostname = COALESCE(host.name, winlog.computer_name)
  | EVAL username = COALESCE(user.name, winlog.event_data.SubjectUserName)
  | WHERE hostname IS NOT NULL
  // Scheduled task indicators
  | WHERE (
      cmd RLIKE ".*(schtasks.*/create|schtasks.*/change|at .*[0-9]{1,2}:[0-9]{2}).*" OR
      cmd RLIKE ".*(crontab -e|crontab -l|crontab -).*" OR
      filepath RLIKE ".*(\\\\Tasks\\\\|/etc/cron|/var/spool/cron).*"
    )
  // Aggregate by host
  | STATS 
      task_ops = COUNT(*),
      users = COUNT_DISTINCT(username),
      sample_commands = VALUES(SUBSTRING(cmd, 0, 150))
    BY hostname
  // Threshold: 2+ task operations
  | WHERE task_ops >= 2
  // Risk scoring
  | EVAL risk_score = CASE(
      task_ops >= 5, 85,
      task_ops >= 3, 75,
      65
    )
  | KEEP hostname, task_ops, users, risk_score, sample_commands
enabled: false
tags:
  - Persistence
  - Scheduled Tasks
  - Cron
  - Privilege Escalation
severity: medium
risk_score: 70
references:
  - "https://attack.mitre.org/techniques/T1053/"
  - "https://attack.mitre.org/techniques/T1053/005/"
nist: ["AU-2", "CM-3", "SI-4"]
gdpr: ["Article 32"]
pci_dss: ["10.2.5", "10.2.7"]
mitre_attack:
  TA0003:
    T1053: ["T1053.005", "T1053.003"]