# converted_sigma/linux/auditd/path/lnx_auditd_systemd_service_creation_esql.yml
rule_id: e63550eb-f80d-4acf-890d-571bef23a55b
name: "Systemd Service Creation Detection"
description: |
  Detects creation of systemd services which could be used by adversaries to execute malicious code.
  Monitors systemd service directories for new .service files.
  Adversaries may create or modify systemd services to establish persistence.
type: esql
params: '{}'
schedule_interval: 10m
index:
  - auditbeat-*
query_language: esql
query: |
  FROM auditbeat-*
  /* auditd PATH type events with CREATE */
  | WHERE event.action == "created" OR (event.dataset RLIKE "(?i)auditd" AND event.type == "creation")
  /* selection: systemd service directories */
  | WHERE (
      file.path RLIKE "(?i)^/usr/lib/systemd/system/"
      OR file.path RLIKE "(?i)^/etc/systemd/system/"
      OR file.path RLIKE "(?i)/\\.config/systemd/user/"
    )
  | STATS
      event_count = COUNT(*),
      users = VALUES(user.name),
      processes = VALUES(process.name),
      executables = VALUES(process.executable),
      file_paths = VALUES(file.path),
      file_names = VALUES(file.name),
      hosts = VALUES(host.name),
      raw_samples = VALUES(event.original)
    BY host.name, user.name, bucket = DATE_TRUNC(20 minutes, @timestamp)
  | WHERE event_count >= 1
  | KEEP agent_id, bucket, host.name, user.name, event_count, users, processes, executables, file_paths, file_names, hosts, raw_samples
  | SORT event_count DESC
enabled: false
tags:
  - attack.privilege-escalation
  - attack.persistence
  - attack.t1543.002
  - detection.threat-hunting
severity: medium
risk_score: 60
references:
  - "https://github.com/redcanaryco/atomic-red-team/blob/f339e7da7d05f6057fdfcdd3742bfcf365fee2a9/atomics/T1543.002/T1543.002.md"
  - "https://attack.mitre.org/techniques/T1543/002/"
nist: ["SI-4", "AU-12", "CM-7"]
pci_dss: ["10.2", "2.2"]
gdpr: ["32.1.b"]
mitre_attack:
  TA0004:
    T1543: ["T1543.002"]
  TA0003:
    T1543: ["T1543.002"]
