# converted_sigma/windows/file/file-event/Scheduled Task Created - FileCreation.yml
rule_id: 575f30db-ce59-4b81-95c8-ded888037a26
name: "Scheduled Task Created - FileCreation"
description: |
  Detects creation of Scheduled Tasks via file creation in Windows task folders:
  \Windows\System32\Tasks\, \Windows\SysWOW64\Tasks\, or \Windows\Tasks\.
  Attackers commonly create scheduled tasks for execution, persistence, and privilege escalation.
  False positives: Normal Windows/admin behavior. Use as hunting baseline and correlate with task content,
  creator process, and unusual timing/user context.
type: esql
params: '{}'
schedule_interval: 10m
index:
  - logs-winlog*
query_language: esql
query: |
  FROM logs-winlog*
  // File event telemetry
  | WHERE event.category == "file" OR event.dataset RLIKE "(?i).*file.*"
  // Evaluate path to simplify logic
  | EVAL p = COALESCE(winlog.event_data.OriginalFileName, "")
  // Selection: TargetFilename contains task folders
  | WHERE
      p LIKE "*:\\\\Windows\\\\System32\\\\Tasks\\\\*"
      OR p LIKE "*:\\\\Windows\\\\SysWOW64\\\\Tasks\\\\*"
      OR p LIKE "*:\\\\Windows\\\\Tasks\\\\*"
  | STATS
      event_count = COUNT(*),
      users = VALUES(winlog.event_data.TargetUserName),
      processes = VALUES(winlog.event_data.ProcessName),
      executables = VALUES(winlog.event_data.ProcessName),
      parents = VALUES(winlog.event_data.CallerProcessName),
      parent_executables = VALUES(winlog.event_data.CallerProcessName),
      file_paths = VALUES(winlog.event_data.OriginalFileName),
      cmdlines = VALUES(COALESCE(winlog.event_data.CommandLine, winlog.event_data.CommandLine, ""))
    BY 
      host.name, 
      bucket = DATE_TRUNC(20m, @timestamp)
  | WHERE event_count >= 1
  | KEEP 
      agent_id, bucket, 
      host.name, 
      event_count, 
      users, 
      processes, 
      executables, 
      parents, 
      parent_executables, 
      cmdlines, 
      file_paths
  | SORT event_count DESC
enabled: false
tags:
  - attack.execution
  - attack.persistence
  - attack.privilege-escalation
  - attack.t1053.005
  - attack.s0111
  - car.2013-08-001
  - detection.threat-hunting
severity: low
risk_score: 20
references:
  - "https://center-for-threat-informed-defense.github.io/summiting-the-pyramid/analytics/task_scheduling/"
  - "https://posts.specterops.io/abstracting-scheduled-tasks-3b6451f6a1c5"
nist: ["SI-4", "AU-12", "CM-6"]
pci_dss: ["10.2"]
gdpr: ["32.1.b"]
mitre_attack:
  TA0003:
    T1053: ["T1053.005"]
